# platform/macos/com.weewx.weewxd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.weewx.weewxd</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- CRITICAL: restart if it exits -->
    <key>KeepAlive</key>
    <true/>

	<key>ProgramArguments</key>
	<array>
    	<string>/usr/local/weewx-venv/bin/weewxd</string>
    	<string>/usr/local/etc/weewx/weewx-data/weewx.conf</string> 
	</array>

    <!-- Capture logs -->
<key>StandardOutPath</key>
<string>/Library/Logs/weewx.out</string>

<key>StandardErrorPath</key>
<string>/Library/Logs/weewx.err</string>

    <!-- Give it a sane environment -->
<key>EnvironmentVariables</key>
<dict>
    <key>PATH</key>
    <string>/usr/local/weewx-venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>WEEWX_BINDIR</key>
    <string>/usr/local/weewx-venv/bin</string>
    <key>SSL_CERT_FILE</key>
    <string>/usr/local/weewx-venv/lib/python3.14/site-packages/certifi/cacert.pem</string>
</dict>

</dict>
</plist>
